# crowdsec/scripts/systemd.d/crowdsec-fastly-bouncer.service
[Unit]
Description=CrowdSec Fastly Bouncer
After=crowdsec.service
Requires=crowdsec.service

[Service]
WorkingDirectory=/app/cs/bin
Type=exec
Environment=LC_ALL=C LANG=C
ExecStartPre=/app/cs/bin/crowdsec-fastly-bouncer -c /app/cs/etc/crowdsec/bouncers/crowdsec-fastly-bouncer.yaml -t
ExecStart=/app/cs/bin/crowdsec-fastly-bouncer -c /app/cs/etc/crowdsec/bouncers/crowdsec-fastly-bouncer.yaml
Restart=always
RestartSec=60
User=%i

[Install]
WantedBy=multi-user.target