# crowdsec/scripts/systemd.d/http2udpRelay.service
[Unit]
Description=HTTP to UDP Forwarder

[Service]
ExecStart=/app/cs/usr/local/bin/http2udpRelay
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
