# crowdsec/scripts/systemd.d/rsyslog.service
[Unit]
Description=rsyslog

[Service]
ExecStart=/app/.global/nix-env/root/bin/rsyslogd -f /app/cs/etc/rsyslog.conf
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target